# OrangePi/babycam-push-optimized.service
[Unit]
Description=BabyCam RTSP Push (Software/Optimized)
After=network.target mediamtx.service

[Service]
ExecStartPre=/bin/bash -lc 'for i in {1..60}; do ss -lnt | grep -q ":8554 " && exit 0; sleep 1; done; echo "RTSP 8554 not ready"; exit 1'

ExecStart=/usr/bin/ffmpeg \
  -f v4l2 \
  -input_format mjpeg \
  -framerate 10 \
  -video_size 640x480 \
  -i /dev/video0 \
  -c:v libx264 \
  -preset ultrafast \
  -tune zerolatency \
  -profile:v baseline \
  -pix_fmt yuv420p \
  -g 10 \
  -keyint_min 10 \
  -b:v 500k \
  -maxrate 600k \
  -bufsize 300k \
  -f rtsp \
  -rtsp_transport tcp \
  rtsp://127.0.0.1:8554/babycam

Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
